# rocky10/docker/storage/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
stdout_logfile=/var/log/sshd.out
stderr_logfile=/var/log/sshd.err

[program:storage-init]
command=/usr/local/bin/storage-init.sh
autostart=true
# storage-init execs ganesha.nfsd in the foreground, so this program is the
# long-lived NFS server. Restart it if the server dies (but not when it exits 0
# via the DISABLE_NFS_AUTOSTART idle path).
autorestart=unexpected
exitcodes=0
startsecs=5
stdout_logfile=/var/log/storage.out
stderr_logfile=/var/log/storage.err
